# deploy/oam-hub.service
[Unit]
Description=OAM The Hub Gateway
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fizobia
EnvironmentFile=/opt/fizobia/.env.server
ExecStart=/usr/bin/python3 -m app.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
